# packaging/deb/caddy-nithronos/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create necessary directories
        mkdir -p /etc/caddy/Caddyfile.d
        mkdir -p /var/log/caddy
        
        # Set proper permissions
        chown -R caddy:caddy /var/log/caddy
        chmod 750 /var/log/caddy
        
        # Install default Caddyfile if not exists
        if [ ! -f /etc/caddy/Caddyfile ]; then
            cp /usr/share/nithronos/caddy/Caddyfile /etc/caddy/Caddyfile
            chown root:caddy /etc/caddy/Caddyfile
            chmod 644 /etc/caddy/Caddyfile
        fi
        
        # Validate configuration
        if caddy validate --config /etc/caddy/Caddyfile >/dev/null 2>&1; then
            echo "Caddy configuration validated successfully"
        else
            echo "Warning: Caddy configuration validation failed, service may not start"
        fi
        
        # Enable and reload Caddy
        if systemctl is-enabled caddy >/dev/null 2>&1; then
            systemctl reload caddy || true
        else
            systemctl enable caddy || true
            systemctl start caddy || true
        fi
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
