# deploy/systemd/nithronos-firstboot-banner.service
[Unit]
Description=NithronOS First-boot Console Banner
After=nosd.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'ip=$(ip -4 route get 1.1.1.1 2>/dev/null | awk '\''/src/ {for(i=1;i<=NF;i++) if ($i=="src") {print $(i+1); exit}}'\'' || true); [ -z "$ip" ] && ip=$(hostname -I | awk '{print $1}'); echo -e "\nNithronOS UI → https://$ip/  (self-signed)\nUse the one-time OTP shown above to enter Setup.\n"'
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target


