# deploy/systemd/nos-announce-otp.service
[Unit]
Description=Announce first-boot OTP to all ttys
After=getty.target multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'otp=$(cat /run/nos/firstboot-otp 2>/dev/null || true); \
  [ -n "$otp" ] && msg="NithronOS first-boot OTP: $otp (valid 15m)"; \
  [ -n "$msg" ] && { echo "$msg" | tee /dev/console | wall -n; \
                    printf "%s\n" "$msg" > /etc/motd.d/10-nithronos-otp; } || true'
RemainAfterExit=no


