# deploy/systemd/nos-announce-otp.path
[Unit]
Description=Watch for first-boot OTP and announce to consoles

[Path]
PathChanged=/run/nos/firstboot-otp

[Install]
WantedBy=multi-user.target


