# deploy/systemd/nos-firstboot-otp.service
[Unit]
Description=NithronOS First-Boot OTP Display
After=getty@tty1.service multi-user.target network-online.target nosd.service
Wants=multi-user.target
ConditionPathExists=/var/lib/nos/firstboot.json

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/nos-firstboot-otp.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
