# packaging/deb/nos-web/debian/postinst
#!/bin/sh
set -e

install -d -o root -g root -m 0755 /etc/caddy
install -d -o root -g root -m 0755 /etc/systemd/system/caddy.service.d
install -d -o root -g root -m 0755 /usr/share/nithronos/web

# Own web assets by www-data (best effort)
chown -R www-data:www-data /usr/share/nithronos/web 2>/dev/null || true

# Enable and start caddy
systemctl daemon-reload || true
systemctl enable nos-tls-gen.service caddy.service || true
systemctl restart nos-tls-gen.service || true
systemctl restart caddy.service || true
systemctl preset nos-tls-gen.service caddy.service || true

#DEBHELPER#

exit 0


