# packaging/deb/nos-web/debian/postinst
#!/bin/sh
set -e

# Ensure directories
install -d -o root -g root -m 0755 /etc/caddy
install -d -o root -g root -m 0755 /etc/systemd/system/caddy.service.d
install -d -o root -g root -m 0755 /usr/share/nithronos/web
install -d -o root -g root -m 0755 /var/log/caddy

# Own web assets by www-data (best effort)
chown -R www-data:www-data /usr/share/nithronos/web 2>/dev/null || true

# Set ownership for caddy log directory
chown caddy:caddy /var/log/caddy 2>/dev/null || true

# Detect if we're running in chroot (no systemd)
CHROOT=0
if [ ! -d /run/systemd/system ]; then
    CHROOT=1
    echo "Running in chroot environment, skipping Caddy validation and service operations"
fi

# Create systemd override for reliable ordering
cat >/etc/systemd/system/caddy.service.d/override.conf <<'EOF'
[Unit]
Wants=network-online.target
After=network-online.target nosd.service

[Service]
ExecStartPre=/usr/bin/caddy validate --config /etc/caddy/Caddyfile
Restart=always
RestartSec=2s
EOF

# Only validate and enable Caddy if not in chroot
if [ "$CHROOT" -eq 0 ]; then
    # Validate and enable Caddy
    /usr/bin/caddy validate --config /etc/caddy/Caddyfile || {
        echo "Warning: Caddy validation failed, but continuing installation"
        echo "This may be expected if TLS certificates are not yet generated"
    }
    systemctl daemon-reload || true
    systemctl enable --now caddy || true
else
    echo "Skipping Caddy validation in chroot environment"
fi

#DEBHELPER#

exit 0