# packaging/deb/nos-web/debian/postinst
#!/bin/sh
set -e

install -d -o root -g root -m 0755 /etc/caddy
# Ensure our include exists at /etc/caddy/Caddyfile.nithronos (installed by deb)
if [ ! -f /etc/caddy/Caddyfile ]; then
  printf 'import /etc/caddy/Caddyfile.nithronos\n' > /etc/caddy/Caddyfile
elif ! grep -q 'Caddyfile\.nithronos' /etc/caddy/Caddyfile; then
  printf '\nimport /etc/caddy/Caddyfile.nithronos\n' >> /etc/caddy/Caddyfile
fi

# Own web assets by www-data (best effort)
chown -R www-data:www-data /usr/share/nithronos/web 2>/dev/null || true

# Enable and start caddy
systemctl daemon-reload || true
systemctl enable --now caddy || true

#DEBHELPER#

exit 0


