# packaging/deb/nithronos/nos-btrfs-scrub@.service
[Unit]
Description=Btrfs scrub for %i

[Service]
Type=oneshot
ExecStart=/usr/bin/btrfs scrub start -B %i


