# packaging/deb/nithronos/nos-smart-scan.timer
[Unit]
Description=Weekly SMART health scan timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target


