# packaging/deb/nithronos/nos-smart-scan.service
[Unit]
Description=Weekly SMART health scan

[Service]
Type=oneshot
ExecStart=/usr/lib/nithronos/smart-scan.sh


