# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nithronos-installer-welcome.service
[Unit]
Description=NithronOS Installer Welcome Message
After=multi-user.target
ConditionKernelCommandLine=nos.install=1

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/nithronos/installer-welcome

[Install]
WantedBy=multi-user.target
