# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nithronos-issue-gen.service
[Unit]
Description=Generate NithronOS issue file with system info
After=network-online.target
Wants=network-online.target
Before=getty@.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/nithronos/generate-issue

[Install]
WantedBy=multi-user.target
