# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nithronos-otp-gen.service
[Unit]
Description=NithronOS OTP Generator for First Boot
Before=nosd.service
After=local-fs.target
ConditionPathExists=!/etc/nos/otp

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/nithronos/generate-otp
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
