# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nithronos-firstboot-banner.service
[Unit]
Description=NithronOS First-boot Console Banner
After=nosd.service

[Service]
Type=oneshot
ExecStart=/usr/lib/nithronos/show-firstboot-banner
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
