# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nithronos-ensure-services.service
[Unit]
Description=NithronOS Ensure Services Running
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nithronos-ensure-services.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
