# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nos-agent.service
[Unit]
Description=NithronOS Privileged Agent (nos-agent)
After=network.target

[Service]
User=root
Group=root
ExecStart=/usr/sbin/nos-agent
Restart=on-failure
RestartSec=2s
StandardOutput=journal+console
StandardError=journal+console
# Simplified for ISO environment - full sandboxing in production

[Install]
WantedBy=multi-user.target
