# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nos-agent.service
[Unit]
Description=NithronOS Privileged Agent (nos-agent)
After=network.target
ConditionPathExists=/usr/sbin/nos-agent

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/sbin/nos-agent
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
# Simplified for ISO environment - full sandboxing in production

[Install]
WantedBy=multi-user.target
