# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nos-agent.service
[Unit]
Description=NithronOS Privileged Agent (nos-agent)
After=network.target
Before=nosd.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/sbin/nos-agent
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Ensure socket is accessible
ExecStartPost=/bin/sh -c 'sleep 1; chmod 666 /run/nos-agent.sock 2>/dev/null || true'

[Install]
WantedBy=multi-user.target