# packaging/iso/debian/config/includes.chroot/etc/systemd/system/nosd.service
[Unit]
Description=NithronOS API Server (nosd)
After=network.target nos-agent.service
Wants=nos-agent.service

[Service]
User=nos
Group=nos
ExecStart=/usr/bin/nosd
Restart=on-failure
RestartSec=2s
StandardOutput=journal+console
StandardError=journal+console
# Simplified for ISO environment - full sandboxing in production

[Install]
WantedBy=multi-user.target
