# packaging/iso/debian/config/includes.chroot/etc/systemd/system/otp-announce.path
[Unit]
Description=NithronOS: re-announce OTP on change
After=nosd.service

[Path]
PathChanged=/var/lib/nos/state/firstboot.json
Unit=otp-announce.service

[Install]
WantedBy=multi-user.target
