# packaging/iso/debian/config/includes.chroot/etc/systemd/system/otp-announce.service
[Unit]
Description=NithronOS: OTP Generation Service
After=multi-user.target
ConditionPathExists=/usr/lib/nithronos/otp-notify
# OTP is now shown in web UI, not on console
# To re-enable terminal display, create /etc/nos/enable-terminal-otp
ConditionPathExists=/etc/nos/enable-terminal-otp

[Service]
Type=oneshot
ExecStart=/usr/lib/nithronos/otp-notify
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
