# packaging/iso/debian/config/includes.chroot/etc/systemd/system/otp-announce.service
[Unit]
Description=NithronOS: announce first-boot OTP to console
After=multi-user.target getty@tty1.service
ConditionPathExists=/usr/lib/nithronos/otp-notify
# Run after system and console are ready

[Service]
Type=oneshot
# Wait longer to ensure terminal is fully ready and user can see
ExecStartPre=/bin/sleep 5
ExecStart=/usr/lib/nithronos/otp-notify
# Also display OTP multiple times with pauses for visibility
ExecStartPost=/bin/sleep 2
ExecStartPost=/usr/lib/nithronos/otp-notify
ExecStartPost=/bin/sleep 2
ExecStartPost=/usr/lib/nithronos/otp-notify
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes
# Try multiple times if OTP not ready yet
Restart=on-failure
RestartSec=2s
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
