# packaging/deb/nosd/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Ensure ownership of configuration and data directories
        if getent passwd nos >/dev/null 2>&1; then
            chown -R nos:nos /etc/nos /var/lib/nos || true
        fi
        # Create runtime and state/config directories via systemd helpers
        if command -v systemd-sysusers >/dev/null 2>&1; then
            systemd-sysusers || true
        fi
        if command -v systemd-tmpfiles >/dev/null 2>&1; then
            systemd-tmpfiles --create || true
        fi
        ;;
esac

systemctl daemon-reload || true
systemctl enable nosd.service || true
systemctl restart nosd.service || true

exit 0


