# packaging/deb/nosd/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Ensure ownership of configuration and data directories
        chown -R nosd:nosd /etc/nos /var/lib/nos || true
        # Create runtime and state/config directories via systemd helpers
        systemd-sysusers || true
        systemd-tmpfiles --create || true
        ;;
esac

systemctl daemon-reload || true
systemctl enable nosd.service || true
systemctl restart nosd.service || true

exit 0


