# packaging/deb/nosd/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create system user and directories for state/config
        systemd-sysusers || true
        systemd-tmpfiles --create || true
        chown -R nos:nos /etc/nos /var/lib/nos 2>/dev/null || true
        chmod 0770 /etc/nos /var/lib/nos 2>/dev/null || true
        [ -f /etc/nos/secret.key ] && chmod 600 /etc/nos/secret.key || true
        
        # Enable OTP announcement service
        systemctl daemon-reload || true
        systemctl enable otp-announce.service || true
        ;;
esac

#DEBHELPER#

exit 0


