# packaging/deb/nos-agent/debian/postinst
#!/bin/sh
set -e

systemctl daemon-reload || true
systemctl enable nos-agent.service || true
systemctl restart nos-agent.service || true

exit 0


