# deploy/firstboot/print-url@.service
[Unit]
Description=Print first-boot URL and OTP for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/bash -c '
  set -euo pipefail;
  mkdir -p /etc/nos;
  OTP=$(head -c 12 /dev/urandom | base64 | tr -dc A-Za-z0-9 | head -c 8);
  IP=$(hostname -I | awk "{print \$1}");
  echo "First-run URL: https://$IP/";
  echo "One-time setup code: $OTP";
  echo "$OTP" > /etc/nos/first-run;
'
RemainAfterExit=true

[Install]
WantedBy=multi-user.target


